# docker/cron/crontab
* * * * * cd /var/www && php artisan schedule:run >> /dev/stdout 2>&1
* * * * * echo $(date) >> /dev/stdout 2>&1
